# pickle-03/crontab/root
# m h  dom mon dow   command
*/15 *  * * * /root/useradd.sh
