# abtech-tracker-02/crontab/root
# m h  dom mon dow   command
30 04 * * * /root/mysqldumper
