# abtech-tracker-02/crontab/root
# m h  dom mon dow   command
30 04 * * * /root/mysqldumper
55 14 * * 6 /root/dump2afs
